# disable_lid_wake/disable_lid_wake.service
[Unit]
Description=Disable lid wake
 
[Service]
ExecStart={{INSTALL_DIR}}/disable_lid_wake/disable_lid_wake.sh
 
[Install]
WantedBy=multi-user.target