# disable_lid_wake/disable_lid_wake.service
[Unit]
Description=Disable lid wake
 
[Service]
ExecStart=/opt/ar18/disable_lid_wake/disable_lid_wake.sh
 
[Install]
WantedBy=multi-user.target